# scripts/linux/systemd/qbitxd.service
# Q-BitX Node Service (standalone SHA-256d with post-quantum address support)
# This is a template - actual paths are set during installation
#
# Variables to replace:
#   {{INSTALL_DIR}} - Installation directory (default: /spiralpool)
#   {{QBX_DIR}} - Q-BitX binary directory
#   {{POOL_USER}} - Pool system user (default: spiraluser)

[Unit]
Description=Q-BitX Daemon
Documentation=https://qbitx.org
After=network-online.target
Wants=network-online.target
# Prevent restart loops: max 5 restarts in 60 minutes (nodes can take 30+ min to start)
StartLimitIntervalSec=3600
StartLimitBurst=5

[Service]
Type=forking
User={{POOL_USER}}
Group={{POOL_USER}}
ExecStart={{QBX_DIR}}/bin/qbitxd -daemon -conf={{INSTALL_DIR}}/qbx/qbitx.conf -datadir={{INSTALL_DIR}}/qbx
ExecStop={{QBX_DIR}}/bin/qbitx-cli -conf={{INSTALL_DIR}}/qbx/qbitx.conf stop
PIDFile={{INSTALL_DIR}}/qbx/qbitxd.pid
Restart=always
RestartSec=30
TimeoutStartSec=infinity
TimeoutStopSec=600
LimitNOFILE=65535

# NOTE: Systemd security hardening intentionally omitted.
# Some blockchain daemons crash with SIGABRT under modern systemd
# hardening options (PrivateTmp, ProtectSystem, etc.)

[Install]
WantedBy=multi-user.target
