# scripts/linux/systemd/spiraldash.service
# Spiral Pool Dashboard Service
# This is a template - actual paths are set during installation
#
# Variables to replace:
#   {{INSTALL_DIR}} - Installation directory (default: /spiralpool)
#   {{POOL_USER}} - Pool system user (default: spiraluser)
#   {{DASHBOARD_PORT}} - Dashboard port (default: 1618)
#   {{METRICS_TOKEN}} - Prometheus metrics auth token

[Unit]
Description=Spiral Pool Dashboard
Documentation=https://github.com/SpiralPool/Spiral-Pool
After=network-online.target
# R-11 FIX: Increased from 300s to 600s (match sentinel).
# Dashboard depends on stratum, which can take 5+ minutes after cold boot.
# 300s was too aggressive — dashboard hit restart limit before stratum was ready.
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User={{POOL_USER}}
Group={{POOL_USER}}
WorkingDirectory={{INSTALL_DIR}}/dashboard

# Single worker required: scan progress uses in-memory state that cannot be shared across processes
# 4 threads handle concurrent requests (polls, page loads, API calls) within the single process
# HTTPS: enable from dashboard Management tab (runs enable-https.sh to add TLS args)
ExecStart={{INSTALL_DIR}}/dashboard/venv/bin/gunicorn --bind 0.0.0.0:{{DASHBOARD_PORT}} --workers 1 --threads 4 --timeout 120 dashboard:app
# POST-BOOT HEALTH CHECK: After reboot, gunicorn worker can deadlock during fork
# (post-fork threading lock inheritance). The master process stays alive and systemd
# sees the service as "active", but no HTTP requests are served. This loop polls
# the dashboard for up to 60 seconds; if it never responds, kills the main process
# so Restart=always can recover it. Uses -k (insecure) for HTTPS self-signed certs.
ExecStartPost=/bin/bash -c '\
  proto=http; \
  grep -q "\\-\\-certfile" /proc/$MAINPID/cmdline 2>/dev/null && proto=https; \
  for i in $(seq 1 30); do \
    sleep 2; \
    curl -sfk -o /dev/null "$${proto}://127.0.0.1:{{DASHBOARD_PORT}}/" 2>/dev/null && exit 0; \
  done; \
  echo "spiraldash: health check failed after 60s — killing hung worker" >&2; \
  kill -KILL $MAINPID 2>/dev/null; \
  exit 1'

# === INSTALLATION PATHS ===
# Tell dashboard where to find data files (miners.json, etc.)
Environment="SPIRALPOOL_INSTALL_DIR={{INSTALL_DIR}}"

# === POOL API & METRICS ===
# Metrics auth token for Prometheus /metrics endpoint on Spiral Stratum
Environment="SPIRAL_METRICS_TOKEN={{METRICS_TOKEN}}"

# === DASHBOARD AUTHENTICATION ===
# Dashboard requires authentication by default (set password on first access)
Environment="DASHBOARD_AUTH_ENABLED=true"
Environment="DASHBOARD_SESSION_LIFETIME=24"
# Optional: Set API key for programmatic access (uncomment and set your own key)
# Environment="DASHBOARD_API_KEY=your-secure-api-key-here"
# Optional: Set allowed CORS origins (empty = same-origin only, most secure)
# Environment="DASHBOARD_CORS_ORIGINS="

# === SELF-HEALING / AUTO-RESTART ===
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# === SECURITY HARDENING ===
# Filesystem protection
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
# /etc/systemd/system: enable-https.sh (sudo) patches spiraldash.service via sed -i.
# ProtectSystem=strict makes /etc read-only in this namespace; without this exception
# sed can't create its temp file ("Read-only file system"). File permissions (root:root 644)
# still prevent the unprivileged dashboard process from writing directly.
ReadWritePaths={{INSTALL_DIR}}/dashboard {{INSTALL_DIR}}/config {{INSTALL_DIR}}/data {{INSTALL_DIR}}/certs /home/{{POOL_USER}}/.spiralpool /etc/systemd/system /var/lib/apt /var/cache/apt /var/log/apt

# Privilege restrictions
# NoNewPrivileges must be 'no' — dashboard uses sudo systemctl restart/stop
# for service control via web UI (sudoers rules in /etc/sudoers.d/spiralpool-dashboard)
NoNewPrivileges=no
# Bounding set: only the capabilities sudo needs to escalate (setuid/setgid/audit).
# Empty CapabilityBoundingSet= would block sudo entirely ("unable to change to root gid").
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_DAC_OVERRIDE CAP_AUDIT_WRITE CAP_FOWNER
AmbientCapabilities=

# Kernel protection
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectProc=invisible

# Network restrictions (AF_INET6 kept for dual-stack compat; IPv6 disabled at OS level)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# System call filtering
# NOTE: SystemCallFilter is deliberately OMITTED here. Using SystemCallFilter
# implicitly forces NoNewPrivileges=yes (systemd requirement for seccomp filters).
# This dashboard needs NoNewPrivileges=no because it uses sudo for service control
# via web UI (sudoers rules in /etc/sudoers.d/spiralpool-dashboard).
SystemCallArchitectures=native

# Memory protection
# NOTE: MemoryDenyWriteExecute omitted — Python/gunicorn may use mmap with W+X
# for C extension loading. Safe because NoNewPrivileges=no already implies reduced
# memory protection guarantees.
LockPersonality=yes

# Misc hardening
RestrictRealtime=yes
RestrictSUIDSGID=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
