# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:spiralstratum]
command=/spiralpool/bin/spiralstratum -config /spiralpool/config/config.yaml
directory=/spiralpool
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/spiraluser"

# NOTE: Dashboard and Sentinel run as separate containers in docker-compose.yml.
# They are NOT managed by supervisord to avoid duplicate processes and port conflicts.
# For standalone single-container use, start them manually or add [program:] sections.
